# files/files_asus-ax53u/etc/uci-defaults/93_wan-setup
uci -q batch << EOI
set network.wan='interface'
set network.wan.proto='dhcp'
set network.wan.peerdns=0
add_list network.wan.dns="74.82.42.42"
add_list network.wan.dns="8.8.8.8"
add_list network.wan.dns="2001:470:20::2"
add_list network.wan.dns="2001:4860:4860::8888"
EOI
