# files/files_asus-ax53u/etc/uci-defaults/95_ntp-enable-server
uci -q batch << EOI
set system.ntp.enable_server='1'
set system.ntp.interface='lan'
EOI
