# files/etc/uci-defaults/90_led-mapping-wps-to-vpn
CHECK=`uci show system | grep "@led" | grep "name='vpn'" | wc -l`

if [ "$CHECK" -gt 0 ]; then
    exit 0
fi

uci -q batch << EOI
add system led
set system.@led[-1].name='vpn'
set system.@led[-1].sysfs='green:wps'
set system.@led[-1].trigger='netdev'
set system.@led[-1].dev='tun0'
set system.@led[-1].mode='link'
EOI
