# files/etc/uci-defaults/99_led-mapping-wps-to-vpn
uci -q batch << EOI
add system led
set system.@led[-1].name='vpn'
set system.@led[-1].sysfs='green:wps'
set system.@led[-1].trigger='netdev'
set system.@led[-1].dev='tun0'
set system.@led[-1].mode='link'
EOI
