# files/etc/uci-defaults/99_ntp-enable-server
uci -q batch << EOI
set system.ntp.enable_server='1'
EOI
